# supervisord.conf
[supervisord]
nodaemon=true
logfile=/opt/supervisord.log

[program:gdalcubes_server]
command=/opt/gdalcubes/build/src/gdalcubes_server
autorestart=unexpected
exitcodes=0
startsecs=0
